# third/logstash-forwarder/supervisor_logstash-forwarder.conf
[program:logstash-forwarder]
directory=/
command=/etc/init.d/logstash-forwarder start
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log